# data/repos/phabricator-repos.path
[Unit]
Description=Auto-update phabricator repos in cgit

[Path]
PathChanged=/opt/phacility/git/

[Install]
WantedBy=default.target
