# data/repos/phabricator-repos.service
[Unit]
Description=Update wiki with repository information

[Service]
Type=oneshot
ExecStart=/usr/bin/phabricator-tools-update-repos

