# data/repos/phabricator-repo-wiki.timer
[Unit]
Description=Update wiki of repos

[Timer]
Unit=phabricator-repo-wiki.service
OnCalendar=*-*-* *:00:00


[Install]
WantedBy=timers.target
