# data/tar-repos/phabricator-tar-repos.service
[Unit]
Description=Nightly tar build of phabricator repos

[Service]
Type=oneshot
ExecStart=/usr/bin/phabricator-tools-tar-gz-repos
