# data/reports/load-reports.path
[Unit]
Description=Auto-update phab reports

[Path]
PathChanged=/opt/phab/reports/

[Install]
WantedBy=default.target
