# data/reports/load-reports.service
[Unit]
Description=Load inbox reports into phabricator

[Service]
Type=oneshot
ExecStart=/usr/bin/phabricator-tools-load-reports
