# utilities/phab-utilities.service
[Unit]
Description=Run phabricator utility processes

[Service]
Type=oneshot
ExecStart=/usr/bin/phab-utilities
