# data/backup/backup.timer
[Unit]
Description=Phabricator Content Backup

[Timer]
Unit=phabricator-backup.service
OnCalendar=*-*-* 03:00:00

[Install]
WantedBy=timers.target
