# data/backup/phabricator-backup.service
[Unit]
Description=Perform scheduled backup of phabricator

[Service]
Type=oneshot
ExecStart=/usr/bin/phabricator-tools-backup
