# data/cleansing/phab-cleanse.service
[Unit]
Description=Run phabricator cleanse processes

[Service]
Type=oneshot
ExecStart=/usr/bin/phab-cleanse
