# utilities/system/systemd-units/phab-weekly.service
[Unit]
Description=Phab Weekly Service

[Service]
Type=oneshot
ExecStart=/opt/phabricator-tools/utilities/api/scripts/weekly.sh
