# utilities/system/systemd-units/phab-weekly.timer
[Unit]
Description=Phab Weekly Tasks

[Timer]
Unit=phab-weekly.service
OnCalendar=Mon *-*-* 00:00:00

[Install]
WantedBy=timers.target
