# utilities/system/systemd-units/phab-daily.service
[Unit]
Description=Phab Daily Service

[Service]
Type=oneshot
ExecStart=/opt/phabricator-tools/utilities/system/scripts/timed.sh daily
